# jellyfin_playlists_sync.timer
[Unit]
Description=Sync playlists automatically

[Timer]
OnCalendar=23:45:00

[Install]
WantedBy=timers.target
